# docs/remote-codex-runner/codex-watchdog.service
[Unit]
Description=Ensure Codex tmux session exists

[Service]
Type=oneshot
ExecStart=%h/bin/start-eom-loop.sh
